# config/includes.chroot_after_packages/etc/systemd/system/automount.service
[Unit]
Description=automount known partitions on boot

[Service]
Type=oneshot
ExecStart=/usr/sbin/automount
RemainAfterExit=true

[Install]
WantedBy=multi-user.target

